# deploy/supervisor.conf
# /etc/supervisor/conf.d/frappe-bench.conf
# Supervisor configuration for Frappe Bench

[group:frappe-bench-web]
programs=frappe-bench-frappe-web,frappe-bench-redis-cache,frappe-bench-redis-queue,frappe-bench-redis-socketio

[program:frappe-bench-frappe-web]
command=/home/frappe/frappe-bench/env/bin/gunicorn \
    -b 127.0.0.1:8000 \
    -w 4 \
    --timeout 120 \
    --max-requests 2000 \
    --max-requests-jitter 200 \
    frappe.app:application
directory=/home/frappe/frappe-bench/sites
user=frappe
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/frappe/frappe-bench/logs/web.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
startsecs=5

[program:frappe-bench-frappe-worker-default]
command=/home/frappe/frappe-bench/env/bin/python -m frappe.utils.bench_helper frappe worker --queue default
directory=/home/frappe/frappe-bench
user=frappe
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/frappe/frappe-bench/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:frappe-bench-frappe-schedule]
command=/home/frappe/frappe-bench/env/bin/python -m frappe.utils.bench_helper frappe schedule
directory=/home/frappe/frappe-bench
user=frappe
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/frappe/frappe-bench/logs/schedule.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:frappe-bench-redis-cache]
command=redis-server /home/frappe/frappe-bench/config/redis_cache.conf
priority=1
user=frappe
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/redis-cache.log

[program:frappe-bench-redis-queue]
command=redis-server /home/frappe/frappe-bench/config/redis_queue.conf
priority=1
user=frappe
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/redis-queue.log
